# packaging/vibepaper.service
[Unit]
Description=vibepaper Wayland wallpaper daemon
Documentation=https://github.com/maxischmaxi/vibepaper
PartOf=graphical-session.target
After=graphical-session.target
# Start inside any graphical session: a Wayland one (wlr/GNOME/KDE backends need
# WAYLAND_DISPLAY) or an X11 one (x11/GNOME/KDE backends need DISPLAY). The `|`
# makes these triggering conditions, so the unit starts if EITHER is set.
ConditionEnvironment=|WAYLAND_DISPLAY
ConditionEnvironment=|DISPLAY

[Service]
Type=simple
ExecStart=/usr/bin/vibepaper daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
